# lingxi-demo.service
[Unit]
Description=灵犀销转AI演示应用
Documentation=https://github.com/yourusername/lingxi_demo
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=forking
User=nginx
Group=nginx
WorkingDirectory=/var/www/lingxi_demo
ExecStartPre=/usr/sbin/nginx -t
ExecStart=/usr/sbin/nginx
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
KillSignal=SIGQUIT
TimeoutStopSec=5
KillMode=mixed
PrivateTmp=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target 